# scripts/can0-setup.service
# =============================================================================
# can0-setup.service — Bring up physical CAN interface at boot (Raspberry Pi)
#
# Install:
#   sudo cp scripts/can0-setup.service /etc/systemd/system/
#   sudo systemctl enable can0-setup
#   sudo systemctl start can0-setup
#
# Requires USB-CAN adapter (CANable/gs_usb) plugged in at boot.
# For MCP2515 HAT, use /boot/config.txt dtoverlay instead.
# =============================================================================

[Unit]
Description=Bring up CAN0 interface at 500 kbps
After=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/sbin/modprobe gs_usb
ExecStart=/bin/sh -c "ip link set can0 type can bitrate 500000 && ip link set can0 up"
ExecStop=/sbin/ip link set can0 down

[Install]
WantedBy=multi-user.target
